# magento2/ubuntu/crontab
*/1 * * * * www-data /usr/local/bin/php /var/www/dev/magento2/bin/magento cron:run
*/1 * * * * www-data /usr/local/bin/php /var/www/dev/magento2/bin/magento indexer:reindex
